# app/shared/uci-defaults/01-expand-rootfs
#!/bin/sh
# /etc/uci-defaults/01-expand-rootfs
# Auto-expand rootfs partition and overlay filesystem on first boot.
#
# Automatically works for any device naming scheme (/dev/mmcblk*, /dev/sd*,
# /dev/nvme*, /dev/vd*) without hardcoding device paths or breaking disk
# suffixes (like /dev/mmcblk1 -> /dev/mmcblk).
#
# Removes any subsequent partitions (e.g. p3, p4) after the root/overlay
# partition (p2) so that p2 can be expanded to occupy the whole disk.

logger -t "expand-rootfs" "Starting auto-expansion of rootfs..."

# ---------------------------------------------------------------------------
# 1. Find the physical (non-loop) block device that backs the root or overlay
#    We check in order: /overlay -> /rom -> / (excluding any /dev/loop* devices).
# ---------------------------------------------------------------------------
TARGET_DEV="$(block info 2>/dev/null | grep 'MOUNT="/overlay"' | grep -v '^/dev/loop' | head -n 1 | cut -d: -f1)"

if [ -z "$TARGET_DEV" ]; then
    TARGET_DEV="$(block info 2>/dev/null | grep 'MOUNT="/rom"' | grep -v '^/dev/loop' | head -n 1 | cut -d: -f1)"
fi

if [ -z "$TARGET_DEV" ]; then
    TARGET_DEV="$(block info 2>/dev/null | grep 'MOUNT="/"' | grep -v '^/dev/loop' | head -n 1 | cut -d: -f1)"
fi

if [ -z "$TARGET_DEV" ]; then
    logger -t "expand-rootfs" "Error: Could not find target root/overlay block device."
    exit 1
fi

# ---------------------------------------------------------------------------
# 2. Safely extract DISK and PART without corrupting numbers in disk names
# ---------------------------------------------------------------------------
PART="$(echo "$TARGET_DEV" | grep -oE '[0-9]+$')"

if [ -z "$PART" ]; then
    logger -t "expand-rootfs" "Error: Could not determine partition number from $TARGET_DEV."
    exit 1
fi

DISK="${TARGET_DEV%p"$PART"}"
if [ "$DISK" = "$TARGET_DEV" ]; then
    DISK="${TARGET_DEV%"$PART"}"
fi

if [ -z "$DISK" ] || [ "$DISK" = "$TARGET_DEV" ]; then
    logger -t "expand-rootfs" "Error: Could not determine base disk device from $TARGET_DEV."
    exit 1
fi

logger -t "expand-rootfs" "Base disk: $DISK, Target Partition: $PART ($TARGET_DEV)"

# ---------------------------------------------------------------------------
# 3. Wipe any partitions AFTER the target partition (e.g. remove p3, p4 after p2)
#    So that partition $PART can be expanded to occupy the remaining whole disk.
# ---------------------------------------------------------------------------
logger -t "expand-rootfs" "Checking for partitions after partition $PART on $DISK..."

parted -s "$DISK" unit s print 2>/dev/null \
    | awk -v p="$PART" '$1 ~ /^[0-9]+$/ && $1 > p {print $1}' \
    | while read -r pnum; do
        logger -t "expand-rootfs" "Removing partition $pnum after partition $PART to allow full disk expansion..."
        
        # Determine canonical device path for the partition being removed
        if [ -e "${DISK}p${pnum}" ]; then
            dev="${DISK}p${pnum}"
        elif [ -e "${DISK}${pnum}" ]; then
            dev="${DISK}${pnum}"
        else
            dev=""
        fi
        
        if [ -n "$dev" ]; then
            logger -t "expand-rootfs" "Unmounting $dev..."
            umount -l "$dev" 2>/dev/null || true
            swapoff "$dev" 2>/dev/null || true
        fi

        # Remove partition from disk
        parted -s "$DISK" rm "$pnum" 2>/dev/null || true
        
        # Explicitly ask kernel to drop partition device node
        if [ -n "$dev" ] && command -v partx >/dev/null 2>&1; then
            partx -d "$dev" 2>/dev/null || true
        fi
      done

# Notify kernel that subsequent partitions were removed
partprobe "$DISK" 2>/dev/null || partx -u "$DISK" 2>/dev/null || true
sleep 1

# ---------------------------------------------------------------------------
# 4. Fix misplaced GPT backup header
# ---------------------------------------------------------------------------
logger -t "expand-rootfs" "Fixing GPT backup header if necessary..."
echo "Fix" | parted ---pretend-input-tty "$DISK" print >/dev/null 2>&1 || true

partprobe "$DISK" 2>/dev/null || partx -u "$DISK" 2>/dev/null || true
sleep 1

# ---------------------------------------------------------------------------
# 5. Expand partition $PART to fill 100% of available space on the disk
# ---------------------------------------------------------------------------
logger -t "expand-rootfs" "Resizing partition $PART to 100%..."
parted -s "$DISK" resizepart "$PART" 100%

partprobe "$DISK" 2>/dev/null || partx -u "$DISK" 2>/dev/null || block auto
sleep 1

# ---------------------------------------------------------------------------
# 6. Resize the filesystem inside the expanded partition
# ---------------------------------------------------------------------------
logger -t "expand-rootfs" "Determining filesystem to resize..."

RESIZE_DEV=""

# Find loop device backing /overlay (standard OpenWrt pattern)
OVERLAY_LOOP="$(awk '$2 == "/overlay" && $1 ~ /^\/dev\/loop/ {print $1}' /proc/mounts 2>/dev/null | head -n 1)"

if [ -n "$OVERLAY_LOOP" ]; then
    RESIZE_DEV="$OVERLAY_LOOP"
    logger -t "expand-rootfs" "Found loop-mounted overlay on $RESIZE_DEV inside expanded partition $TARGET_DEV"
    # Refresh loop device capacity so it sees the enlarged underlying partition p2
    losetup -c "$RESIZE_DEV" 2>/dev/null || true
else
    # Otherwise (e.g., ext4 root partition), resize the target partition directly
    RESIZE_DEV="$TARGET_DEV"
    logger -t "expand-rootfs" "Resizing physical block device $RESIZE_DEV directly..."
fi

# Detect filesystem type of the device to resize
RESIZE_FSTYPE="$(awk -v dev="$RESIZE_DEV" '$1 == dev {print $3}' /proc/mounts 2>/dev/null | head -n 1)"
if [ -z "$RESIZE_FSTYPE" ] || [ "$RESIZE_FSTYPE" = "rootfs" ]; then
    RESIZE_FSTYPE="$(block info "$RESIZE_DEV" 2>/dev/null | grep -oE 'TYPE="[^"]+"' | cut -d'"' -f2)"
fi

logger -t "expand-rootfs" "Resizing filesystem on $RESIZE_DEV (type: '${RESIZE_FSTYPE:-unknown}')..."

case "$RESIZE_FSTYPE" in
    f2fs)
        if command -v resize.f2fs >/dev/null 2>&1; then
            resize.f2fs "$RESIZE_DEV"
        else
            logger -t "expand-rootfs" "Error: resize.f2fs not found. Cannot resize f2fs filesystem."
        fi
        ;;
    ext4|ext3|ext2)
        if command -v resize2fs >/dev/null 2>&1; then
            resize2fs -f "$RESIZE_DEV"
        else
            logger -t "expand-rootfs" "Error: resize2fs not found. Cannot resize ext filesystem."
        fi
        ;;
    *)
        logger -t "expand-rootfs" "Warning: unknown filesystem type '$RESIZE_FSTYPE'. Attempting resize2fs and resize.f2fs..."
        resize2fs -f "$RESIZE_DEV" 2>/dev/null || resize.f2fs "$RESIZE_DEV" 2>/dev/null || true
        ;;
esac

logger -t "expand-rootfs" "Rootfs expansion completed successfully."

# (Append custom script logic below)
