# app/shared/uci-defaults/01-expand-rootfs
#!/bin/sh
# /etc/uci-defaults/99-expand-rootfs
# Auto-expand rootfs partition and filesystem on first boot

logger -t "expand-rootfs" "Starting auto-expansion of rootfs..."

# Find the root block device
ROOT_DEV="$(block info | grep -e 'MOUNT="/rom"' -e 'MOUNT="/"' | head -n 1 | cut -d: -f1)"

if [ -z "$ROOT_DEV" ]; then
    logger -t "expand-rootfs" "Error: Could not find root device via block info."
    exit 1
fi

# Extract disk and partition number (handles /dev/mmcblk1p2 -> disk=/dev/mmcblk1, part=2)
DISK="$(echo "$ROOT_DEV" | sed -E 's/p[0-9]+$//')"
PART="$(echo "$ROOT_DEV" | grep -oE '[0-9]+$')"

if [ -z "$DISK" ] || [ -z "$PART" ]; then
    logger -t "expand-rootfs" "Error: Could not parse disk/partition."
    exit 1
fi

logger -t "expand-rootfs" "Fixing GPT backup header if necessary..."
# When a smaller image is flashed to a larger SD card, the GPT backup header
# ends up in the middle of the disk. 'parted resizepart 100%' then fails with
# "Can't have overlapping partitions" because 100% is computed relative to the
# misplaced backup header, not the real disk end.
#
# Fix strategy using ONLY parted (default lightweight tool):
#   1. Use 'parted ---pretend-input-tty' with 'Fix' to relocate the backup GPT
#      header to the real end of the disk. This is a safe no-op on MBR disks.
#   2. Force the kernel to re-read the updated partition table so parted sees
#      the corrected layout before we call resizepart.
#   3. Use 'parted unit s print' to compute the true last free sector and pass
#      an explicit sector number '<N>s' to resizepart instead of '100%'.
echo "Fix" | parted ---pretend-input-tty "$DISK" print >/dev/null 2>&1 || true

# Force kernel to see the updated partition table
partprobe "$DISK" 2>/dev/null || partx -u "$DISK" 2>/dev/null || true
sleep 1  # give udev a moment to settle

logger -t "expand-rootfs" "Expanding partition $PART on disk $DISK..."
# To avoid parted miscomputing '100%' on disks where the GPT backup header was just relocated,
# we compute the last available free sector explicitly using sector units.
PART_END="$(parted -s "$DISK" unit s print 2>/dev/null | awk -v p="$PART" '$1 == p {print $3}' | grep -oE '[0-9]+')"
FREE_END="$(parted -s "$DISK" unit s print free 2>/dev/null | grep -i 'free space' | tail -n 1 | awk '{print $2}' | grep -oE '[0-9]+')"

if [ -n "$FREE_END" ] && [ -n "$PART_END" ] && [ "$FREE_END" -gt "$PART_END" ] 2>/dev/null; then
    logger -t "expand-rootfs" "Resizing partition $PART from ${PART_END}s to ${FREE_END}s..."
    parted -s "$DISK" unit s resizepart "$PART" "${FREE_END}s"
else
    logger -t "expand-rootfs" "Fallback: Resizing partition $PART to 100%..."
    parted -s "$DISK" resizepart "$PART" 100%
fi

# Reload partition table
partx -u "$DISK" 2>/dev/null || block auto

logger -t "expand-rootfs" "Resizing filesystem..."
# OpenWrt squashfs overlay is typically attached via loop0
if losetup -a 2>/dev/null | grep -q "/dev/loop0"; then
    losetup -c /dev/loop0
    resize2fs -f /dev/loop0
else
    resize2fs -f "$ROOT_DEV"
fi

logger -t "expand-rootfs" "Rootfs expansion completed successfully."

# (Append custom script logic below)
